# .config/systemd/user/2in1screen.service
k[Unit]
Description=Auto rotation of screen

[Service]
ExecStart=/home/baroud/CS/Project/2in1screen/2in1screen

[Install]
WantedBy=default.target