# .config/systemd/user/2in1screen.service
[Unit]
Description=Auto rotation of screen

[Service]
ExecStart=/home/baroud/bin/2in1screen

[Install]
WantedBy=default.target